# primary/config/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:sing-box]
command=/usr/local/bin/sing-box run -D /var/lib/sing-box -c /etc/sing-box/config.json
environment=ENABLE_DEPRECATED_SPECIAL_OUTBOUNDS="true",ENABLE_DEPRECATED_TUN_ADDRESS_X="true",ENABLE_DEPRECATED_WIREGUARD_OUTBOUND="true"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/sing-box.log
stderr_logfile=/var/log/supervisor/sing-box.err.log
priority=10

[program:fastapi]
command=/usr/bin/uvicorn server.main:app --host 0.0.0.0 --port 9091 --app-dir /app
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/fastapi.log
stderr_logfile=/var/log/supervisor/fastapi.err.log
priority=20

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.err.log
priority=30

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd.err.log
priority=40
